# chat2.0/server/iems5722_app.conf
[program:iems5722_gunicorn3]
command = gunicorn3 -b localhost:8000 -w 4 iems5722_a4:app
directory = /home/ubuntu
user = ubuntu
autostart = true
autorestart = true
stdout_logfile = /home/ubuntu/app.log
redirect_stderr = true

[program:iems5722_celery]
command = celery -A task.celery worker
directory = /home/ubuntu
user = ubuntu
autostart = true
autorestart = true
stdout_logfile = /home/ubuntu/app.log
redirect_stderr = true

[program:iems5722_socketio]
command = python3 iems5722_socketio.py
directory = /home/ubuntu
user = ubuntu
autostart = true
autorestart = true
stdout_logfile = /home/ubuntu/flask-socketio_app.log
redirect_stderr = true
